# 7.2-nginx/etc/supervisor.d/php-fpm.ini
[program:phpfpm]
command=php-fpm7 --fpm-config /etc/php7/php-fpm.conf -c /etc/php7
autostart=true
autorestart=true
autorstart=true
startsecs=5
startretries=100
redirect_stderr=true
